# corosync/PLIST
@comment $NetBSD$
bin/corosync-blackbox
etc/corosync/corosync.conf.example
etc/corosync/corosync.conf.example.udpu
etc/init.d/corosync
etc/init.d/corosync-notifyd
etc/logrotate.d/corosync
include/corosync/cfg.h
include/corosync/cmap.h
include/corosync/corodefs.h
include/corosync/corotypes.h
include/corosync/cpg.h
include/corosync/hdb.h
include/corosync/quorum.h
include/corosync/sam.h
include/corosync/totem/totem.h
include/corosync/totem/totemip.h
include/corosync/totem/totempg.h
include/corosync/votequorum.h
lib/libcfg.la
lib/libcmap.la
lib/libcorosync_common.la
lib/libcpg.la
lib/libquorum.la
lib/libsam.la
lib/libtotem_pg.la
lib/libvotequorum.la
lib/pkgconfig/corosync.pc
lib/pkgconfig/libcfg.pc
lib/pkgconfig/libcmap.pc
lib/pkgconfig/libcorosync_common.pc
lib/pkgconfig/libcpg.pc
lib/pkgconfig/libquorum.pc
lib/pkgconfig/libsam.pc
lib/pkgconfig/libtotem_pg.pc
lib/pkgconfig/libvotequorum.pc
man/man3/cmap_context_get.3
man/man3/cmap_context_set.3
man/man3/cmap_dec.3
man/man3/cmap_delete.3
man/man3/cmap_dispatch.3
man/man3/cmap_fd_get.3
man/man3/cmap_finalize.3
man/man3/cmap_get.3
man/man3/cmap_inc.3
man/man3/cmap_initialize.3
man/man3/cmap_iter_finalize.3
man/man3/cmap_iter_init.3
man/man3/cmap_iter_next.3
man/man3/cmap_set.3
man/man3/cmap_track_add.3
man/man3/cmap_track_delete.3
man/man3/cpg_context_get.3
man/man3/cpg_context_set.3
man/man3/cpg_dispatch.3
man/man3/cpg_fd_get.3
man/man3/cpg_finalize.3
man/man3/cpg_initialize.3
man/man3/cpg_iteration_finalize.3
man/man3/cpg_iteration_initialize.3
man/man3/cpg_iteration_next.3
man/man3/cpg_join.3
man/man3/cpg_leave.3
man/man3/cpg_local_get.3
man/man3/cpg_mcast_joined.3
man/man3/cpg_membership_get.3
man/man3/cpg_model_initialize.3
man/man3/cpg_zcb_alloc.3
man/man3/cpg_zcb_free.3
man/man3/cpg_zcb_mcast_joined.3
man/man3/quorum_context_get.3
man/man3/quorum_context_set.3
man/man3/quorum_dispatch.3
man/man3/quorum_fd_get.3
man/man3/quorum_finalize.3
man/man3/quorum_getquorate.3
man/man3/quorum_initialize.3
man/man3/quorum_trackstart.3
man/man3/quorum_trackstop.3
man/man3/sam_data_getsize.3
man/man3/sam_data_restore.3
man/man3/sam_data_store.3
man/man3/sam_finalize.3
man/man3/sam_hc_callback_register.3
man/man3/sam_hc_send.3
man/man3/sam_initialize.3
man/man3/sam_mark_failed.3
man/man3/sam_register.3
man/man3/sam_start.3
man/man3/sam_stop.3
man/man3/sam_warn_signal_set.3
man/man3/votequorum_context_get.3
man/man3/votequorum_context_set.3
man/man3/votequorum_dispatch.3
man/man3/votequorum_fd_get.3
man/man3/votequorum_finalize.3
man/man3/votequorum_getinfo.3
man/man3/votequorum_initialize.3
man/man3/votequorum_qdevice_master_wins.3
man/man3/votequorum_qdevice_poll.3
man/man3/votequorum_qdevice_register.3
man/man3/votequorum_qdevice_unregister.3
man/man3/votequorum_qdevice_update.3
man/man3/votequorum_setexpected.3
man/man3/votequorum_setvotes.3
man/man3/votequorum_trackstart.3
man/man3/votequorum_trackstop.3
man/man5/corosync.conf.5
man/man5/votequorum.5
man/man8/cmap_keys.8
man/man8/cmap_overview.8
man/man8/corosync.8
man/man8/corosync-blackbox.8
man/man8/corosync-cfgtool.8
man/man8/corosync-cmapctl.8
man/man8/corosync-cpgtool.8
man/man8/corosync-keygen.8
man/man8/corosync-notifyd.8
man/man8/corosync_overview.8
man/man8/corosync-quorumtool.8
man/man8/cpg_overview.8
man/man8/quorum_overview.8
man/man8/sam_overview.8
man/man8/votequorum_overview.8
sbin/corosync
sbin/corosync-cfgtool
sbin/corosync-cmapctl
sbin/corosync-cpgtool
sbin/corosync-keygen
sbin/corosync-notifyd
sbin/corosync-quorumtool
share/doc/corosync/AUTHORS
share/doc/corosync/html/cmap_context_get.3.html
share/doc/corosync/html/cmap_context_set.3.html
share/doc/corosync/html/cmap_dec.3.html
share/doc/corosync/html/cmap_delete.3.html
share/doc/corosync/html/cmap_dispatch.3.html
share/doc/corosync/html/cmap_fd_get.3.html
share/doc/corosync/html/cmap_finalize.3.html
share/doc/corosync/html/cmap_get.3.html
share/doc/corosync/html/cmap_inc.3.html
share/doc/corosync/html/cmap_initialize.3.html
share/doc/corosync/html/cmap_iter_finalize.3.html
share/doc/corosync/html/cmap_iter_init.3.html
share/doc/corosync/html/cmap_iter_next.3.html
share/doc/corosync/html/cmap_keys.8.html
share/doc/corosync/html/cmap_overview.8.html
share/doc/corosync/html/cmap_set.3.html
share/doc/corosync/html/cmap_track_add.3.html
share/doc/corosync/html/cmap_track_delete.3.html
share/doc/corosync/html/corosync.8.html
share/doc/corosync/html/corosync-blackbox.8.html
share/doc/corosync/html/corosync-cfgtool.8.html
share/doc/corosync/html/corosync-cmapctl.8.html
share/doc/corosync/html/corosync.conf.5.html
share/doc/corosync/html/corosync-cpgtool.8.html
share/doc/corosync/html/corosync-keygen.8.html
share/doc/corosync/html/corosync-notifyd.8.html
share/doc/corosync/html/corosync_overview.8.html
share/doc/corosync/html/corosync-quorumtool.8.html
share/doc/corosync/html/cpg_context_get.3.html
share/doc/corosync/html/cpg_context_set.3.html
share/doc/corosync/html/cpg_dispatch.3.html
share/doc/corosync/html/cpg_fd_get.3.html
share/doc/corosync/html/cpg_finalize.3.html
share/doc/corosync/html/cpg_initialize.3.html
share/doc/corosync/html/cpg_iteration_finalize.3.html
share/doc/corosync/html/cpg_iteration_initialize.3.html
share/doc/corosync/html/cpg_iteration_next.3.html
share/doc/corosync/html/cpg_join.3.html
share/doc/corosync/html/cpg_leave.3.html
share/doc/corosync/html/cpg_local_get.3.html
share/doc/corosync/html/cpg_mcast_joined.3.html
share/doc/corosync/html/cpg_membership_get.3.html
share/doc/corosync/html/cpg_model_initialize.3.html
share/doc/corosync/html/cpg_overview.8.html
share/doc/corosync/html/cpg_zcb_alloc.3.html
share/doc/corosync/html/cpg_zcb_free.3.html
share/doc/corosync/html/cpg_zcb_mcast_joined.3.html
share/doc/corosync/html/index.html
share/doc/corosync/html/quorum_context_get.3.html
share/doc/corosync/html/quorum_context_set.3.html
share/doc/corosync/html/quorum_dispatch.3.html
share/doc/corosync/html/quorum_fd_get.3.html
share/doc/corosync/html/quorum_finalize.3.html
share/doc/corosync/html/quorum_getquorate.3.html
share/doc/corosync/html/quorum_initialize.3.html
share/doc/corosync/html/quorum_overview.8.html
share/doc/corosync/html/quorum_trackstart.3.html
share/doc/corosync/html/quorum_trackstop.3.html
share/doc/corosync/html/sam_data_getsize.3.html
share/doc/corosync/html/sam_data_restore.3.html
share/doc/corosync/html/sam_data_store.3.html
share/doc/corosync/html/sam_finalize.3.html
share/doc/corosync/html/sam_hc_callback_register.3.html
share/doc/corosync/html/sam_hc_send.3.html
share/doc/corosync/html/sam_initialize.3.html
share/doc/corosync/html/sam_mark_failed.3.html
share/doc/corosync/html/sam_overview.8.html
share/doc/corosync/html/sam_register.3.html
share/doc/corosync/html/sam_start.3.html
share/doc/corosync/html/sam_stop.3.html
share/doc/corosync/html/sam_warn_signal_set.3.html
share/doc/corosync/html/votequorum.5.html
share/doc/corosync/html/votequorum_context_get.3.html
share/doc/corosync/html/votequorum_context_set.3.html
share/doc/corosync/html/votequorum_dispatch.3.html
share/doc/corosync/html/votequorum_fd_get.3.html
share/doc/corosync/html/votequorum_finalize.3.html
share/doc/corosync/html/votequorum_getinfo.3.html
share/doc/corosync/html/votequorum_initialize.3.html
share/doc/corosync/html/votequorum_overview.8.html
share/doc/corosync/html/votequorum_qdevice_master_wins.3.html
share/doc/corosync/html/votequorum_qdevice_poll.3.html
share/doc/corosync/html/votequorum_qdevice_register.3.html
share/doc/corosync/html/votequorum_qdevice_unregister.3.html
share/doc/corosync/html/votequorum_qdevice_update.3.html
share/doc/corosync/html/votequorum_setexpected.3.html
share/doc/corosync/html/votequorum_setvotes.3.html
share/doc/corosync/html/votequorum_trackstart.3.html
share/doc/corosync/html/votequorum_trackstop.3.html
share/doc/corosync/INSTALL
share/doc/corosync/LICENSE
share/doc/corosync/README.recovery
share/doc/corosync/SECURITY
